# misc/rank.timer
[Unit]
Description=Code Rank Timer

[Timer]
# auto run crawler service
OnUnitInactiveSec=10800seconds
