# misc/rank.timer
[Unit]
Description=Code Rank Timer

[Timer]
OnCalendar=*:0/1
